# roles/system_notifications/templates/system-alerts.service.j2
[Unit]
Description=System Alerts Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/system-health-check
User=root

[Install]
WantedBy=multi-user.target
