# templates/pbs-backup.timer.j2
[Unit]
Description=Run TrueNAS PBS Backup Daily
Requires=pbs-backup.service

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target
